# autorotate.service
[Unit]
Description=Auto-rotate screen and Wacom input based on accelerometer
Documentation=https://github.com/pastelrx/lifebook-autorotate
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
# Edit this path to point at wherever you cloned the repo:
ExecStart=%h/lifebook-autorotate/autorotate.sh
Restart=on-failure
RestartSec=3
# The script needs the X display; these are usually inherited in a user
# session, but set them explicitly if the service can't reach X:
# Environment=DISPLAY=:0

[Install]
WantedBy=graphical-session.target
