# autorotate.service
[Unit]
Description=Auto-rotate screen and Wacom input based on tablet mode + accelerometer
Documentation=https://github.com/pastelrx/lifebook-autorotate
After=graphical-session.target
# Keep retrying until X is ready instead of giving up after a few fast failures.
StartLimitIntervalSec=0

[Service]
Type=simple
# Runs the script via the symlink the README's install step drops into ~/.local/bin.
# (Point this at your clone directly if you prefer, e.g. %h/lifebook-autorotate/autorotate.sh)
ExecStart=%h/.local/bin/autorotate.sh
# X11 access. systemd --user usually imports these from the running session, but
# set them explicitly so the service also works if it starts before that import.
# Change DISPLAY if yours isn't :0.
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
# If the script starts before X is up, xrandr fails and it exits; we just retry.
Restart=on-failure
RestartSec=3

[Install]
# graphical-session.target is often inactive on XFCE, so bind to default.target
# (reached when the user systemd instance starts at login) to actually autostart.
WantedBy=default.target
